# layer-data/vector-db/launchd/com.bebop.pgvector-main.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bebop.pgvector-main</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/postgresql@16/bin/postgres</string>
    <string>-D</string>
    <string>/Users/thestudio/data/memory-main/postgres/data</string>
    <string>-c</string>
    <string>listen_addresses=127.0.0.1</string>
    <string>-c</string>
    <string>port=55432</string>
    <string>-c</string>
    <string>wal_level=replica</string>
    <string>-c</string>
    <string>archive_mode=on</string>
    <string>-c</string>
    <string>archive_command=test ! -f /Users/thestudio/data/memory-main/postgres/wal_archive/%f &amp;&amp; cp %p /Users/thestudio/data/memory-main/postgres/wal_archive/%f</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>LC_ALL</key>
    <string>C.UTF-8</string>
    <key>LANG</key>
    <string>C.UTF-8</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/thestudio/data/memory-main/postgres/data</string>
  <key>StandardOutPath</key>
  <string>/Users/thestudio/data/memory-main/logs/pgvector-main.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thestudio/data/memory-main/logs/pgvector-main.err</string>
  <key>UserName</key>
  <string>thestudio</string>
  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>10</integer>
  <key>LowPriorityIO</key>
  <true/>
  <key>LowPriorityBackgroundIO</key>
  <true/>
</dict>
</plist>
