# layer-inference/optillm-local/systemd/optillm-local.service
[Unit]
Description=OptiLLM Local (Orin)
After=network.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/opt/homelab/optillm-local/layer-inference/optillm-local
EnvironmentFile=/etc/optillm-local/env
ExecStart=/opt/homelab/optillm-local/layer-inference/optillm-local/.venv/bin/optillm \
  --host 0.0.0.0 \
  --port 4040 \
  --model ${OPTILLM_MODEL} \
  --approach none \
  --optillm-api-key ${OPTILLM_API_KEY}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
