# services/vector-db/launchd/com.bebop.memory-api-main.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bebop.memory-api-main</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/thestudio/optillm-proxy/layer-data/vector-db/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>app.main:app</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>55440</string>
    <string>--workers</string>
    <string>1</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MEMORY_DB_HOST</key>
    <string>127.0.0.1</string>
    <key>MEMORY_DB_PORT</key>
    <string>55432</string>
    <key>MEMORY_DB_USER</key>
    <string>memory_app</string>
    <key>MEMORY_DB_NAME</key>
    <string>memory_main</string>
    <key>MEMORY_DB_AUTO_INIT</key>
    <string>false</string>
    <key>MEMORY_BACKEND</key>
    <string>elastic</string>
    <key>MEMORY_ELASTIC_URL</key>
    <string>http://127.0.0.1:9200</string>
    <key>MEMORY_EMBED_PRIMARY</key>
    <string>studio-nomic-embed-text-v1.5</string>
    <key>MEMORY_EMBED_NOMIC_DIM</key>
    <string>768</string>
    <key>MEMORY_EMBED_NOMIC_QUERY_PREFIX</key>
    <string>search_query:</string>
    <key>MEMORY_EMBED_NOMIC_DOCUMENT_PREFIX</key>
    <string>search_document:</string>
    <key>MEMORY_EMBED_DOCUMENT_BATCH_SIZE</key>
    <string>1</string>
    <key>MEMORY_EMBED_FALLBACK</key>
    <string>studio-mxbai-embed-large-v1</string>
    <key>MEMORY_API_WRITE_BEARER_TOKEN_FILE</key>
    <string>/Users/thestudio/data/memory-main/secrets/memory-api-write-token</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/thestudio/optillm-proxy/layer-data/vector-db</string>
  <key>StandardOutPath</key>
  <string>/Users/thestudio/data/memory-main/logs/memory-api-main.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thestudio/data/memory-main/logs/memory-api-main.err</string>
  <key>UserName</key>
  <string>thestudio</string>
  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>10</integer>
  <key>LowPriorityIO</key>
  <true/>
  <key>LowPriorityBackgroundIO</key>
  <true/>
</dict>
</plist>
