# services/vector-db/launchd/com.bebop.kibana-memory-main.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bebop.kibana-memory-main</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/thestudio/optillm-proxy/layer-data/vector-db/runtime/kibana-current/bin/kibana</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>KBN_PATH_CONF</key>
    <string>/Users/thestudio/optillm-proxy/layer-data/vector-db/runtime/kibana-config</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/thestudio/optillm-proxy/layer-data/vector-db/runtime/kibana-current</string>
  <key>StandardOutPath</key>
  <string>/Users/thestudio/data/memory-main/logs/kibana-memory-main.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thestudio/data/memory-main/logs/kibana-memory-main.err</string>
  <key>UserName</key>
  <string>thestudio</string>
  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>10</integer>
  <key>LowPriorityIO</key>
  <true/>
  <key>LowPriorityBackgroundIO</key>
  <true/>
</dict>
</plist>
