# ops/systemd/litellm-orch.service
[Unit]
Description=LiteLLM Proxy (litellm-orch)
After=network.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/home/christopherbailey/homelab-llm/services/litellm-orch
EnvironmentFile=/home/christopherbailey/homelab-llm/services/litellm-orch/config/env.local
Environment=LITELLM_LOG=INFO
ExecStart=/home/christopherbailey/homelab-llm/services/litellm-orch/.venv/bin/litellm --config /home/christopherbailey/homelab-llm/services/litellm-orch/config/router.yaml --host 0.0.0.0 --port 4000
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
