# platform/ops/systemd/litellm-orch.service
[Unit]
Description=LiteLLM Proxy (litellm-orch)
After=network.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/home/christopherbailey/homelab-llm/services/litellm-orch
EnvironmentFile=/home/christopherbailey/homelab-llm/services/litellm-orch/config/env.local
Environment=LITELLM_LOG=INFO
Environment=ENFORCE_PRISMA_MIGRATION_CHECK=true
ExecStart=/home/christopherbailey/homelab-llm/services/litellm-orch/scripts/run-service.sh
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
