# platform/ops/systemd/optillm-proxy.service
[Unit]
Description=OptiLLM Proxy (optillm-proxy)
After=network.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/home/christopherbailey/homelab-llm/layer-gateway/optillm-proxy
EnvironmentFile=/etc/optillm-proxy/env
ExecStart=/home/christopherbailey/homelab-llm/layer-gateway/optillm-proxy/.venv/bin/optillm --host 127.0.0.1 --port 4020 --base-url http://127.0.0.1:4000/v1 --approach proxy --model mlx-gpt-oss-120b-mxfp4-q4 --optillm-api-key dummy
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
