# platform/ops/systemd/orchestration-cockpit-graph.service
[Unit]
Description=orchestration-cockpit LangGraph local Agent Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/home/christopherbailey
EnvironmentFile=/etc/orchestration-cockpit/graph.env
EnvironmentFile=-/etc/orchestration-cockpit/graph.secret.env
NoNewPrivileges=true
PrivateTmp=true
ExecStart=/bin/bash -lc 'exec "$ORCHESTRATION_COCKPIT_REPO_ROOT/services/orchestration-cockpit/scripts/run_langgraph_dev.sh"'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
