# platform/ops/systemd/orchestration-cockpit-ui.service
[Unit]
Description=orchestration-cockpit Agent Chat UI
After=network-online.target orchestration-cockpit-graph.service
Wants=network-online.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/home/christopherbailey
EnvironmentFile=/etc/orchestration-cockpit/ui.env
NoNewPrivileges=true
PrivateTmp=true
ExecStart=/bin/bash -lc 'exec "$ORCHESTRATION_COCKPIT_REPO_ROOT/services/orchestration-cockpit/scripts/run_agent_chat_ui.sh"'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
