# platform/ops/systemd/searxng.service
[Unit]
Description=SearXNG (searxng)
After=network.target

[Service]
Type=simple
User=christopherbailey
WorkingDirectory=/home/christopherbailey/homelab-llm/services/searxng/app
EnvironmentFile=/etc/searxng/env
ExecStart=/home/christopherbailey/homelab-llm/services/searxng/app/.venv/bin/granian searx.webapp:app --interface wsgi --host 127.0.0.1 --port 8888
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
