# docs/archive/optillm-local/launchd/optillm-local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>optillm.local</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/thestudio/optillm-local/.venv/bin/optillm</string>
      <string>--host</string>
      <string>0.0.0.0</string>
      <string>--port</string>
      <string>4040</string>
      <string>--model</string>
      <string>${OPTILLM_LOCAL_MODEL}</string>
      <string>--approach</string>
      <string>none</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>OPTILLM_API_KEY</key>
      <string>${OPTILLM_LOCAL_API_KEY}</string>
      <key>PYTORCH_ENABLE_MPS_FALLBACK</key>
      <string>1</string>
      <key>HF_HOME</key>
      <string>/Users/thestudio/models/hf</string>
      <key>TRANSFORMERS_CACHE</key>
      <string>/Users/thestudio/models/hf/hub</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/thestudio/optillm-local</string>
    <key>StandardOutPath</key>
    <string>/Users/thestudio/Library/Logs/optillm-local.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thestudio/Library/Logs/optillm-local.err</string>
  </dict>
</plist>
